# Addboard/packages/bsp/UFI001C/adbd.service
[Unit]
Description=ADB Daemon
After=sysinit.target

[Service]
Type=simple
ExecStart=/usr/bin/adbd
Restart=on-failure
RestartSec=2s